# docker/s6-rc.d/frontend_dev/run
#!/command/with-contenv sh
set -e

cd /var/www/funding_call/frontend/

[ "${ENVIRONMENT}" = "development" ] && npm run dev
